# palworld.service
# Generic Systemd Service Template

This document provides a template for creating a systemd service file. You can customize it for your specific application.

## Service File Template

Create a new file named `your-service.service` in your systemd directory (typically `/etc/systemd/system/`):

[Unit]
Description=Your Application Description
After=network.target

[Service]
Type=simple
User=youruser         # Replace with your username
ExecStart=/path/to/your/executable/startup/script.sh      # Replace with your full script path
RemainAfterExit=yes
Restart=on-failure
RestartSec=5
StandardOutput=append:/var/log/yourapp.log
StandardError=append:/var/log/yourapp.log

[Install]
WantedBy=multi-user.target
